# config/aussieshopper-site.ini
[program:aussieshopper]
command=/home/aussieshopper/venv/bin/gunicorn --bind localhost:8010 aussieshopper.wsgi:application
process_name=%(program_name)s_%(process_num)s
numprocs=1
directory=/home/aussieshopper/src
stdout_logfile=/home/aussieshopper/logs/api-gunicorn.log
stderr_logfile=/home/aussieshopper/logs/api-gunicorn-error.log
autostart=true
autorestart=true
user=aussieshopper

killasgroup=true
stopasgroup=true